# root/etc/s6-overlay/s6-rc.d/svc-mod-mstream-mod/run
#!/usr/bin/with-contenv bash

# Run index watch in background
/usr/local/bin/music.sh